# porteus/home/guest/.config/autostart/Verify Operation.desktop
[Desktop Entry]
Exec=/usr/local/bin/verify_operation.sh
Name=Verify Operation
Type=Application
Version=1.0
X-Razor-Need-Tray=true
